# deploy/systemd/binance-futures-agent-cap-rescue.service
[Unit]
Description=Manual CAPUSDT hedge rescue monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/binance-futures-agent/app
Environment=PYTHONPATH=/opt/binance-futures-agent/app/src
ExecStart=/usr/bin/python3 /opt/binance-futures-agent/app/scripts/manual_ops/velvet_rescue.py --symbol CAPUSDT --mode uptrend-short-t --execute --interval 5 --profit-trigger 10 --drawdown-readd 8 --cooldown-seconds 20 --max-imbalance-after-reduce 0.25 --reduce-fraction 0.08 --trend-min-book-delta 2.5 --short-probe-fraction 0.08 --max-short-to-long-ratio 1.02 --short-probe-min-exit-profit-pct 0.18 --error-backoff-seconds 600
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
