# deploy/systemd/binance-futures-agent-paper.timer
[Unit]
Description=Run Binance Futures Agent forward-paper evidence collection every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=false
Unit=binance-futures-agent-paper.service

[Install]
WantedBy=timers.target
