# deploy/systemd/binance-futures-agent-pending-limit-watchdog.service
[Unit]
Description=Binance Futures Agent pending-limit fill watchdog
Documentation=file:/opt/binance-futures-agent/app/docs/deployment.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/binance-futures-agent/app
EnvironmentFile=/etc/binance-futures-agent/env
ExecStart=/opt/binance-futures-agent/.venv/bin/python -m bfa.cli ops pending-limit-watchdog --env-file /etc/binance-futures-agent/env --db /opt/binance-futures-agent/data/agent.sqlite --execute
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/binance-futures-agent /etc/binance-futures-agent
